# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:postgres]
# docker-entrypoint.sh runs as root then drops to postgres user itself.
# Don't set user= here — would break the entrypoint's init step.
command=/usr/local/bin/docker-entrypoint.sh postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    POSTGRES_USER="rtm",
    POSTGRES_PASSWORD="rtm",
    POSTGRES_DB="rtm",
    PGDATA="/var/lib/postgresql/data"

[program:rtm-app]
# wait-for-pg loop before starting uvicorn — DB may still be initializing
command=/bin/bash -c "until pg_isready -h 127.0.0.1 -U rtm -d rtm >/dev/null 2>&1; do echo '[rtm-app] waiting for postgres...'; sleep 2; done; echo '[rtm-app] postgres ready, starting uvicorn'; exec python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8001"
directory=/app
user=root
autostart=true
autorestart=true
priority=20
startretries=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
