# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:postgres]
# docker-entrypoint.sh runs as root then drops to postgres user itself.
# Don't set user= here — would break the entrypoint's init step.
# Tuning flags applied at every boot (overrides defaults; ALTER SYSTEM only
# runs on init so -c is more portable across fresh vs migrated DB).
command=/usr/local/bin/docker-entrypoint.sh postgres
    -c shared_buffers=512MB
    -c work_mem=32MB
    -c maintenance_work_mem=128MB
    -c effective_cache_size=1GB
    -c max_connections=200
    -c random_page_cost=1.1
    -c jit=off
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    POSTGRES_USER="rtm",
    POSTGRES_PASSWORD="rtm",
    POSTGRES_DB="rtm",
    PGDATA="/var/lib/postgresql/data"

[program:rtm-app]
# wait-for-pg loop before starting uvicorn — DB may still be initializing
command=/bin/bash -c "until pg_isready -h 127.0.0.1 -U rtm -d rtm >/dev/null 2>&1; do echo '[rtm-app] waiting for postgres...'; sleep 2; done; echo '[rtm-app] postgres ready, starting uvicorn (4 workers)'; exec python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8001 --workers ${UVICORN_WORKERS:-4}"
directory=/app
user=root
autostart=true
autorestart=true
priority=20
startretries=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
